# 96-microsoft-ergo-4k.rules

ENV{ID_VENDOR}=="Microsoft", ENV{ID_MODEL_ID}=="00db", \
	RUN+="keymap $name carpalx-ms4k"